# services/server-monitor-frontend.service
[Unit]
Description=Server Monitor Next.js Frontend
After=network.target server-monitor-api.service server-monitor-ws.service server-monitor-terminal.service

[Service]
Type=simple
User=server-monitor
Group=server-monitor
WorkingDirectory=/opt/server-monitor/frontend-next
Environment="NODE_ENV=production"
Environment="PORT=9081"
Environment="API_PROXY_TARGET=http://localhost:9083"
Environment="NEXT_PUBLIC_MONITORING_WS_URL=ws://localhost:9085"
Environment="NEXT_PUBLIC_TERMINAL_WS_URL=ws://localhost:9084"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-monitor-frontend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/server-monitor/frontend-next/.next
ReadWritePaths=/var/log/server-monitor

[Install]
WantedBy=multi-user.target
